# roles/amanda_server/templates/systemd/DORMANT_amanda.socket
# amanda.socket 
[Unit]
Description=Amanda Activation Socket

[Socket]
ListenStream=10080
Accept=true

[Install]
WantedBy=sockets.target
